# com.ycintelligence.autoscrape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  YC Intelligence — Weekly Auto-Scraper (macOS)

  INSTALL:
    1. Edit the path to auto_scrape.sh below
    2. Copy this file to your LaunchAgents folder:
       cp com.ycintelligence.autoscrape.plist ~/Library/LaunchAgents/
    3. Load it:
       launchctl load ~/Library/LaunchAgents/com.ycintelligence.autoscrape.plist

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.ycintelligence.autoscrape.plist
    rm ~/Library/LaunchAgents/com.ycintelligence.autoscrape.plist

  This runs every Monday at 9:00 AM.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ycintelligence.autoscrape</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/slyons/Documents/YC Intelligence/auto_scrape.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/slyons/Documents/YC Intelligence/launchd_out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/slyons/Documents/YC Intelligence/launchd_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
